# scripts/StenoMics/stenomics.service
[Unit]
Description=StenoMics automation.

[Service]
Type=oneshot
User=nvidia
Group=nvidia
RemainAfterExit=true
StandardOutput=journal
ExecStart=/bin/bash /usr/bin/stenomics-creation.sh --default

[Install]
WantedBy=default.target