# scripts/StenoMics/stenomics.service
[Unit]
Description=StenoMics automation.

[Service]
Type=idle
User=nvidia
Group=nvidia
ExecStart=/bin/bash /home/nvidia/rendezvous/scripts/stenomics-creation.sh --default

[Install]
WantedBy=multi-user.target