# scripts/overclock-jetson/overclock-jetson.service
[Unit]
Description=Overclock jetson.

[Service]
Type=idel
ExecStart=/bin/bash /usr/bin/overclock-jetson.sh

[Install]
WantedBy=default.target